# scripts/simple-wakeup-handler.service
[Unit]
Description=Simple Wakeup Handler Service
After=network.target
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/webdev/glyphy-board/scripts/simple_wakeup_handler.py
User=root
Group=root
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300
TimeoutStopSec=30

# Restart policy (not needed for oneshot, but included for completeness)
Restart=no
RestartSec=0

# Environment variables
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target 